# rootfs/etc/crontabs/caddy
5 * * * * /usr/local/bin/console core:archive --url=http://localhost/ > /dev/null 2>&1
