# rootfs/etc/crontabs/root
5 * * * * /usr/local/bin/console core:archive --url=http://localhost/
