# v3.8/overlay/etc/crontab
@hourly /usr/bin/console core:archive --url=http://localhost:8080/
